# exef/packaging/local-service/exef-local-service.service
[Unit]
Description=ExEF Local KSeF Service
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/exef/exef.env
RuntimeDirectory=exef
ExecStart=/usr/local/bin/exef-local-service
Restart=on-failure

[Install]
WantedBy=multi-user.target
