# exef/packaging/local-service/exef-local-service.service
[Unit]
Description=ExEF Local KSeF Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/exef-local-service
Restart=on-failure
Environment=LOCAL_SERVICE_HOST=127.0.0.1
Environment=LOCAL_SERVICE_PORT=3030

[Install]
WantedBy=multi-user.target
